# docker/supervisor-celerybeat.conf
; ================================
;  celery beat supervisor example
; ================================

[program:celerybeat]
; Set full path to celery program if using virtualenv

command=/var/app-venv/bin/python manage_ai.py celery beat --pidfile=/tmp/celerybeat.pid --loglevel=debug

directory=/var/app/
user=nobody
numprocs=1
stdout_logfile=/var/log/finmars/celerybeat.log
stderr_logfile=/var/log/finmars/celerybeat.log
autostart=true
autorestart=true
startsecs=10

; if rabbitmq is supervised, set its priority higher so it starts first
priority=999
